# deploy/nextjs.service
[Unit]
Description=Sales Analytics — Next.js Frontend
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/sales-analytics/frontend
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/var/www/sales-analytics/frontend/.env.production
ExecStart=/usr/bin/node /var/www/sales-analytics/frontend/.next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
